# crontab
* * * * * /home/dev_jhesed/ADAUSD_bot/always_up_ada.sh
* * * * * /home/dev_jhesed/ETHUSD_bot/always_up_eth.sh
