# ethusd_bot.conf
description "Pycrypto Bot"
author "jhesed"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec python3 /home/dev_jhesed/ETHUSD_bot/pycryptobot.py